# config/register-bot.service
[Unit]
Description=Telegram Finance Tracker Bot
After=network.target

[Service]
Type=simple
User=YOUR_USER
WorkingDirectory=/home/YOUR_USER/register-track-bot
Environment="TELEGRAM_BOT_TOKEN=your_token_here"
ExecStart=/home/YOUR_USER/register-track-bot/.venv/bin/python run_bot.py
Restart=always
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=register-bot

[Install]
WantedBy=multi-user.target
